# supervisor.conf
[program:shadowsocks]
directory=/root/shadowsocks-munager/
command=/usr/bin/python3 /root/shadowsocks-munager/main.py
autorestart=true
user=root
